# build/supervisord-mainuser.conf
[supervisord]
nodaemon=true
user=mainuser
logfile=/home/mainuser/supervisor/log/supervisord.log
pidfile=/home/mainuser/supervisor/supervisord.pid
childlogdir=/home/mainuser/supervisor/log

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;

[unix_http_server]
file=/home/mainuser/supervisor/supervisor.sock
chown=mainuser:mainuser
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;

[program:sshd]
command=/usr/sbin/sshd -Def /home/mainuser/sshd/sshd_config
priority=10

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;

[include]
files = /home/mainuser/supervisor/conf.d/*.conf
